# general/scripts/excludes/ssc335_initramfs.list
/etc/init.d/S01syslogd
/etc/init.d/S10modules
/etc/init.d/S20urandom
/etc/init.d/S21rngd
/etc/init.d/S40mdev
/etc/init.d/S49ntpd
/etc/init.d/S50dropbear
/etc/init.d/S50snmpd
/lib/libatomic*
/lib/libgcc*
/lib/libjson-c.so.2
/bin/getconf
/bin/ipcinfo
/bin/ldd
/sbin/ldconfig
/sbin/firstboot
/sbin/pqtools
/sbin/sysupgrade
/sbin/telegram_bot
