# general/scripts/excludes/ssc335_initramfs.list
/etc/init.d/S01syslogd
/etc/init.d/S10modules
/etc/init.d/S20urandom
/etc/init.d/S21rngd
/etc/init.d/S40mdev
/etc/init.d/S49ntpd
/etc/init.d/S50dropbear
/etc/init.d/S50snmpd
/lib/libatomic*
/lib/libgcc*
/sbin/ldconfig
/usr/bin/getconf
/usr/bin/ipcinfo
/usr/bin/ldd
/usr/lib/libjson-c.so.2
/usr/sbin/firstboot
/usr/sbin/pqtools
/usr/sbin/sysupgrade
/usr/sbin/telegram_bot
